# com.kevingaasch.smart-ai-router.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kevingaasch.smart-ai-router</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/kevingaasch/ProjectHome/smart-ai-router/.venv/bin/python</string>
    <string>-m</string>
    <string>smart_ai_router</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/kevingaasch/ProjectHome/smart-ai-router</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/kevingaasch</string>
    <key>PATH</key>
    <string>/Users/kevingaasch/.local/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/kevingaasch/ProjectHome/smart-ai-router/logs/server.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/kevingaasch/ProjectHome/smart-ai-router/logs/server.err</string>
</dict>
</plist>
